# snapsync.service
[Unit]
Description=SnapSync - SD Card Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/snapsync
ExecStart=/usr/bin/python3 /home/pi/snapsync/snapsync.py start
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=snapsync

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
